# autorun/316.service
[Unit]
Description=This test-edge.sh will run at startup
After=multi-user.target
[Service]
ExecStart=/root/0144-316.sh
[Install]
WantedBy=multi-user.target